# deploy/it-parser.service
[Unit]
Description=IT tenders parser (icetrade.by -> Telegram)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/opt/parserIT
# Когда вынесете токен в .env — раскомментируйте и укажите путь:
# EnvironmentFile=-/opt/parserIT/.env
ExecStart=/usr/bin/python3 /opt/parserIT/it_parser.py
# Если машина «засыпает», oneshot всё равно не разбудит её — нужен внешний cron
# или отключение сна; этот unit гарантирует корректный запуск, когда пробуждение уже произошло.

[Install]
WantedBy=multi-user.target
